# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=/bin/sh -c "/startServer.sh"

[program:client]
command=/bin/sh -c "/sendData.sh"

